# cron/crontab
# min hour day month weekday command
*/1 * * * * /bin/date
*/1 * * * * curl http://tt-web:5000/tick
*/1 * * * * curl http://tt-web-staging:5000/tick