# cron/crontab
# min hour day month weekday command
*/1 * * * * /bin/date
*/1 * * * * curl http://territorytools:5000/tick